# scripts/deploy/ai-moder.service
[Unit]
Description=AI Moderator local API
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=ai-moder
Group=ai-moder
WorkingDirectory=/opt/ai-moder
EnvironmentFile=/opt/ai-moder/.env
ExecStart=/opt/ai-moder/.venv/bin/uvicorn main_api:app --host 127.0.0.1 --port 8000 --no-proxy-headers
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/ai-moder/logs

[Install]
WantedBy=multi-user.target
